# systemd/nginx-ratelimit-ban/nginx-ratelimit-ban.timer
# /etc/systemd/system/nginx-ratelimit-ban.timer

[Unit]
Description=Run nginx rate-limit auto-ban every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
