# systemd/nginx-404-ban/nginx-404-ban.service
# /etc/systemd/system/nginx-404-ban.service

[Unit]
Description=Auto ban high-frequency 404 scanner IPs from nginx logs
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto_ban_404_scanners.sh --window 120 --threshold 12 --ban-seconds 86400
