# systemd/dbbackup/dbbackup.service
# /etc/systemd/system/dbbackup.service

[Unit]
Description=Database Backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/db_backup.sh